# crypto-signal-bot.service
[Unit]
Description=Crypto Signal Bot - Auto-Healing Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Crypto-signal-bot
Environment="PATH=/root/Crypto-signal-bot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONDONTWRITEBYTECODE=1"
Environment="PYTHONUNBUFFERED=1"

# Cleanup before starting - ensure no stale processes
# Try graceful shutdown first, then force kill if needed
ExecStartPre=/bin/sh -c 'pkill -TERM -f "python.*bot.py" || true'
ExecStartPre=/bin/sleep 2
ExecStartPre=/bin/sh -c 'pkill -9 -f "python.*bot.py" || true'
ExecStartPre=/bin/sleep 1

ExecStart=/root/Crypto-signal-bot/venv/bin/python /root/Crypto-signal-bot/bot.py

# Auto-restart настройки
Restart=always
RestartSec=10
StartLimitInterval=0

# Ensure clean shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# Limits
LimitNOFILE=65536
StandardOutput=append:/root/Crypto-signal-bot/bot.log
StandardError=append:/root/Crypto-signal-bot/bot.log

[Install]
WantedBy=multi-user.target
